# supervisor/conf.d/consul-template.conf
[program:consul-template]
command = /opt/app/consul-template/bin/consul-template -consul-addr "consul:8500" -template "/opt/app/nginx/conf/vhost/nginx.conf.ctmpl:/opt/app/nginx/conf/vhost/nginx.conf:/opt/app/nginx/sbin/nginx -s reload"
user = root
autostart=true
autorestart=true
numprocs=1
numprocs_start=100
priority=99
startsecs=2
startretries=5
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
process_name=%(program_name)s
stdout_logfile=/data/logs/supervisor/%(program_name)s.log
stderr_logfile=/data/logs/supervisor/%(program_name)s_err.log

